# installers/linux/volunteer.service
[Unit]
Description=Volunteer Computing Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=volunteer
Group=volunteer
WorkingDirectory=/opt/volunteer-app
Environment="PYTHONUNBUFFERED=1"
Environment="DJANGO_SETTINGS_MODULE=backend.settings"
ExecStart=/opt/volunteer-app/exp-env/bin/python /opt/volunteer-app/volunteer_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=volunteer-daemon

# Limites de sécurité (privilèges minimum)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# data/ contient db.sqlite3, les fichiers temporaires et les données
ReadWritePaths=/opt/volunteer-app/logs /opt/volunteer-app/data /opt/volunteer-app/pending_requests /opt/volunteer-app/.volunteer
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false

# Performance
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
